# deploy/watchdog.service
[Unit]
Description=frontcrm job_queue watchdog
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/frontcrm/backend
EnvironmentFile=/var/www/frontcrm/backend/.env
ExecStart=/var/www/frontcrm/backend/venv/bin/python -m watchdog
Restart=always
RestartSec=10

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/www/frontcrm/backend
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
